# systemd/tempicon.service
[Unit]
Description=TempIcon - CPU Temperature System Tray Monitor
After=graphical-session-pre.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/share/tempicon/temp_monitor.py
Restart=on-failure
RestartSec=10
Environment="DISPLAY=:0"
Environment="XAUTHORITY=%h/.Xauthority"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus"

[Install]
WantedBy=graphical-session.target
